# etc/supervisor/conf.d/init.conf
[program:init]
command=/root/init.sh
autostart = true
autorestart = false
startretries = 0
priority = 1