# etc/supervisor/conf.d/ulordd.conf
[program:ulordd]
command=/UlordChain/src/ulordd
autostart = true
autorestart = true
priority = 1