# etc/supervisor/conf.d/ulord-pool.conf
[program:ulord-pool]
directory=/ulord-node-stratum-pool
command=npm start
autostart = true
autorestart = true
priority = 20
stdout_logfile = /var/log/ulord-pool.log
loglevel=info

